# deploy/portfolio-api.service
[Unit]
Description=Portfolio API for anikait.page
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/portfolio/backend
EnvironmentFile=/etc/portfolio-api.env
ExecStart=/usr/bin/node src/index.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
